# installer/linux/deb/postinst
#!/bin/sh
set -e

if ! id -u hali >/dev/null 2>&1; then
  useradd --system --no-create-home --shell /usr/sbin/nologin hali
fi

install -d -o hali -g hali /var/lib/hali
install -d -o hali -g hali /var/log/hali
install -d -o hali -g hali /run/hali

systemctl daemon-reload
systemctl enable halid.service
systemctl start halid.service

# Register hali:// URL protocol handler (best-effort; headless servers may lack xdg tooling)
if command -v update-desktop-database >/dev/null 2>&1; then
  update-desktop-database /usr/share/applications/ || true
fi
if command -v xdg-mime >/dev/null 2>&1; then
  xdg-mime default hali.desktop x-scheme-handler/hali || true
fi

exit 0